# Ansible-RoboShop-V1/shipping.service
[Unit]
Description=Shipping Service

[Service]
User=roboshop
Environment=CART_ENDPOINT=cart-dev-01.kiranpanchavati.online:8080
Environment=DB_HOST=mysql-dev-01.kiranpanchavati.online
ExecStart=/bin/java -jar /app/shipping.jar
SyslogIdentifier=shipping

[Install]
WantedBy=multi-user.target
